# system_files/usr/lib/systemd/system/nbc-update-download.service
[Unit]
Description=NBC Update Download

[Service]
Type=oneshot
ExecStart=/usr/bin/nbc update --download-only