# system_files/usr/lib/systemd/system/nbc-update-download.service
[Unit]
Description=NBC Update Download

[Service]
ExecStart=nbc update --download-only